# opencarwings/rootfs/etc/services.d/frpc/run
#!/usr/bin/with-contenv bash
exec /usr/local/bin/frpc -c /etc/frp/frpc.toml