# opencarwings/rootfs/etc/services.d/opencarwings/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
source /etc/cont-init.d/00-log-fix.sh

cd /opt/opencarwings || exit 1

bashio::log.info "OpenCarwings add-on service starting..."

# --- 1. Load Shared Environment ---
# shellcheck source=/etc/opencarwings/env.sh
source /etc/opencarwings/env.sh

DAPHNE_PORT=8000  # Nginx handles external ports 8124/8125

# --- 2. Configure CSRF Trusted Origins ---
# We need to explicitly trust the domains + ports to avoid CSRF errors
CSRF_ORIGINS=""

# Use bashio's native list iteration
DOMAIN_COUNT=$(bashio::config 'trusted_domains | length')

for (( i=0; i < DOMAIN_COUNT; i++ )); do
    domain=$(bashio::config "trusted_domains[${i}]")
    bashio::log.info "Adding CSRF trusted domain: ${domain}"
    # Add raw domain (sometimes needed)
    CSRF_ORIGINS+="https://${domain},"
    CSRF_ORIGINS+="http://${domain},"
    # Add with ports
    CSRF_ORIGINS+="https://${domain}:8125,"
    CSRF_ORIGINS+="http://${domain}:8124,"
done

# Pass this to Django settings
export CSRF_TRUSTED_ORIGINS="${CSRF_ORIGINS%,}"
bashio::log.info "Configured CSRF Trusted Origins: $CSRF_TRUSTED_ORIGINS"

# --- 3. Start OpenCarwings Server ---
bashio::log.info "Starting OpenCarwings server with Daphne on port $DAPHNE_PORT..."

exec daphne --access-log /dev/null -b 0.0.0.0 -p $DAPHNE_PORT \
  carwings.asgi:application