# opencarwings/rootfs/etc/services.d/opencarwings/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
source /etc/cont-init.d/00-log-fix.sh

cd /opt/opencarwings || exit 1

bashio::log.info "OpenCarwings add-on service starting..."

# --- 1. Load Shared Environment ---
# shellcheck source=/etc/opencarwings/env.sh
source /etc/opencarwings/env.sh

DAPHNE_PORT=8000  # Nginx handles external ports 8124/8125

# --- 2. Configure CSRF Trusted Origins ---
# We need to explicitly trust the domains + ports to avoid CSRF errors
CSRF_ORIGINS=""
UNIQUE_DOMAINS_CSV=""
VALID_DOMAIN_COUNT=0

for key in http_domain tcu_domain; do
    domain=$(bashio::config "${key}" '' | xargs)

    if [ -z "$domain" ]; then
        continue
    fi

    if [[ ! "$domain" =~ ^[A-Za-z0-9.-]+$ ]]; then
        bashio::log.warning "Ignoring invalid ${key} value for CSRF origins: '${domain}'"
        continue
    fi

    if [[ ",${UNIQUE_DOMAINS_CSV}" == *",${domain},"* ]]; then
        continue
    fi

    UNIQUE_DOMAINS_CSV+="${domain},"
    VALID_DOMAIN_COUNT=$((VALID_DOMAIN_COUNT + 1))

    bashio::log.info "Adding CSRF trusted domain from ${key}: ${domain}"
    CSRF_ORIGINS+="https://${domain},"
    CSRF_ORIGINS+="http://${domain},"
    CSRF_ORIGINS+="https://${domain}:8125,"
    CSRF_ORIGINS+="http://${domain}:8124,"
done

# Pass this to Django settings
export CSRF_TRUSTED_ORIGINS="${CSRF_ORIGINS%,}"
if [ "$VALID_DOMAIN_COUNT" -eq 0 ]; then
    bashio::log.warning "No valid domains configured in http_domain/tcu_domain for CSRF trusted origins."
else
    bashio::log.info "Configured CSRF Trusted Origins: $CSRF_TRUSTED_ORIGINS"
fi

# --- 3. Start OpenCarwings Server ---
bashio::log.info "Starting OpenCarwings server with Daphne on port $DAPHNE_PORT..."

exec daphne --access-log /dev/null -b 0.0.0.0 -p $DAPHNE_PORT \
  carwings.asgi:application
