# opencarwings/rootfs/etc/services.d/opencarwings/run
#!/bin/sh
cd /opt/opencarwings || exit 1

echo "OpenCarwings add-on service starting..."

# --- 1. Load Configuration Variables ---
TCU_PORT=55230

# Set environment variables for Django
export PSQL_DATABASE=carwings
export PSQL_USER=carwings_user
export PSQL_PASSWORD=secure_password
export PSQL_DATABASE_HOST=localhost
export PSQL_DATABASE_PORT=5432
export REDIS_HOST=localhost
export REDIS_PORT=6379
export PYTHONPATH=/opt/opencarwings
export DJANGO_SETTINGS_MODULE=carwings.settings
export DEBUG=true
export SECRET_KEY=django-insecure-default-key-change-in-production

# --- 2. Start OpenCarwings Server ---
echo "Starting OpenCarwings server with Daphne..."

exec daphne -b 0.0.0.0 -p 8124 \
  carwings.asgi:application