# opencarwings/rootfs/etc/services.d/opencarwings/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash

cd /opt/opencarwings || exit 1

bashio::log.info "OpenCarwings add-on service starting..."

# --- 1. Load Configuration Variables ---
TCU_PORT=55230
DAPHNE_PORT=8124

MODE=$(bashio::config 'connection_mode')

# Check if we should use the internal Nginx proxy (Local Mode)
# In local mode, Nginx ALWAYS runs (either with real certs or self-signed),
# so we must listen on the internal port 8000.
if [ "$MODE" = "local" ]; then
    bashio::log.info "Local mode detected. Switching Daphne to internal port 8000 (proxied by Nginx)."
    DAPHNE_PORT=8000
else
    bashio::log.info "Starting Daphne on standard port 8124 (Direct/Public Mode)."
fi

# Set environment variables for Django
export PSQL_DATABASE=carwings
export PSQL_USER=carwings_user
export PSQL_PASSWORD=secure_password
export PSQL_DATABASE_HOST=localhost
export PSQL_DATABASE_PORT=5432
export REDIS_HOST=localhost
export REDIS_PORT=6379
export PYTHONPATH=/opt/opencarwings
export DJANGO_SETTINGS_MODULE=carwings.settings
export DEBUG=true
export SECRET_KEY=django-insecure-default-key-change-in-production

# --- 2. Configure CSRF Trusted Origins ---
# We need to explicitly trust the domains + ports to avoid CSRF errors
CSRF_ORIGINS=""

# Use bashio's native list iteration (jq doesn't work on bashio config output)
DOMAIN_COUNT=$(bashio::config 'trusted_domains | length')

for (( i=0; i < DOMAIN_COUNT; i++ )); do
    domain=$(bashio::config "trusted_domains[${i}]")
    bashio::log.info "Adding CSRF trusted domain: ${domain}"
    # Add raw domain (sometimes needed)
    CSRF_ORIGINS+="https://${domain},"
    CSRF_ORIGINS+="http://${domain},"
    # Add with ports
    CSRF_ORIGINS+="https://${domain}:8125,"
    CSRF_ORIGINS+="http://${domain}:8124,"
done

# Pass this to Django settings
export CSRF_TRUSTED_ORIGINS="${CSRF_ORIGINS%,}"
bashio::log.info "Configured CSRF Trusted Origins: $CSRF_TRUSTED_ORIGINS"

# --- 3. Start OpenCarwings Server ---
bashio::log.info "Starting OpenCarwings server with Daphne on port $DAPHNE_PORT..."

exec daphne -b 0.0.0.0 -p $DAPHNE_PORT \
  carwings.asgi:application