# opencarwings/rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash

bashio::log.info "Starting Nginx HTTPS proxy..."

# Ensure directories exist
mkdir -p /run/nginx
mkdir -p /data/ssl

# SSL Certificate Handling
# Priority 1: User-provided valid certs from Home Assistant
if [ -f /ssl/fullchain.pem ] && [ -f /ssl/privkey.pem ]; then
    bashio::log.info "Valid SSL certificates found in /ssl/. Using them."
    CERT_FILE="/ssl/fullchain.pem"
    KEY_FILE="/ssl/privkey.pem"
else
    # Priority 2: Self-Signed Fallback
    bashio::log.warning "No valid certificates found in /ssl/. Generating self-signed certificate..."
    
    # Generate only if doesn't exist to avoid annoying browser warnings on every restart
    if [ ! -f /data/ssl/selfsigned.pem ]; then
        openssl req -x509 -nodes -days 3650 -newkey rsa:2048 \
            -keyout /data/ssl/selfsigned.key \
            -out /data/ssl/selfsigned.pem \
            -subj "/C=US/ST=State/L=City/O=OpenCarwings/CN=homeassistant.local"
    fi
    
    CERT_FILE="/data/ssl/selfsigned.pem"
    KEY_FILE="/data/ssl/selfsigned.key"
    bashio::log.warning "Self-signed certificate active. Browsers will show a security warning."
fi

# Build trusted host map from add-on config.
# Requests with a Host header not present in http_domain/tcu_domain are dropped.
TRUSTED_HOSTS_MAP=$'    default 0;\n'
VALID_TRUSTED_DOMAIN_COUNT=0
UNIQUE_DOMAINS_CSV=""

for key in http_domain tcu_domain; do
    domain=$(bashio::config "${key}" '' )
    domain=$(echo "$domain" | xargs)

    if [ -z "$domain" ]; then
        continue
    fi

    if [[ "$domain" =~ ^[A-Za-z0-9.-]+$ ]]; then
        if [[ ",${UNIQUE_DOMAINS_CSV}" == *",${domain},"* ]]; then
            continue
        fi
        UNIQUE_DOMAINS_CSV+="${domain},"
        TRUSTED_HOSTS_MAP+=$'    '"${domain}"$' 1;\n'
        VALID_TRUSTED_DOMAIN_COUNT=$((VALID_TRUSTED_DOMAIN_COUNT + 1))
    else
        bashio::log.warning "Ignoring invalid ${key} value for trusted host allowlist: '${domain}'"
    fi
done

if [ "$VALID_TRUSTED_DOMAIN_COUNT" -eq 0 ]; then
    bashio::log.warning "No valid domains configured in http_domain/tcu_domain. Nginx will fail closed and reject all external requests."
else
    bashio::log.info "Loaded ${VALID_TRUSTED_DOMAIN_COUNT} trusted domain(s) from http_domain/tcu_domain for Nginx host allowlist."
fi

# Generate dynamic Nginx configuration
# - Listen on 8124 (HTTP) -> Single-domain compatible mode (car + UI + HA API)
# - Listen on 8125 (HTTPS) -> Optional direct HTTPS access to UI/API
cat > /etc/nginx/http.d/default.conf <<EOF
map \$uri \$is_carwings_path {
    default 0;
    ~^/WARCondelivbas/it-m_gw10/?$ 1;
}

map \$host \$is_trusted_host {
$TRUSTED_HOSTS_MAP
}

map \$http_user_agent \$is_car_ua {
    default 0;
    ~*NISSAN[[:space:]]+CARWINGS 1;
}

log_format ocw 'client=\$remote_addr:\$remote_port '
               'host="\$host" method=\$request_method uri="\$uri" args="\$args" request="\$request" '
               'status=\$status bytes_sent=\$bytes_sent body_bytes_sent=\$body_bytes_sent '
               'request_length=\$request_length request_time=\$request_time '
               'server_port=\$server_port scheme=\$scheme '
               'referer="\$http_referer" user_agent="\$http_user_agent" '
               'xff="\$http_x_forwarded_for" '
               'upstream_addr="\$upstream_addr" upstream_status="\$upstream_status" '
               'upstream_connect_time="\$upstream_connect_time" '
               'upstream_header_time="\$upstream_header_time" '
               'upstream_response_time="\$upstream_response_time" '
               'connection=\$connection connection_requests=\$connection_requests '
               'ssl_protocol="\$ssl_protocol" ssl_cipher="\$ssl_cipher"';

server {
    listen 8124;
    server_name _;
    
    access_log /dev/stdout ocw;
    error_log /dev/stderr info;

    # Only hosts listed in http_domain/tcu_domain are served.
    if (\$is_trusted_host = 0) {
        return 444;
    }

    location = /WARCondelivbas/it-m_gw10 {
        return 301 /WARCondelivbas/it-m_gw10/;
    }

    location = /WARCondelivbas/it-m_gw10/ {
        if (\$is_car_ua = 0) {
            return 403;
        }
        proxy_pass http://127.0.0.1:8000;
        proxy_set_header Host \$host;
        proxy_set_header X-Real-IP \$remote_addr;
        proxy_set_header X-Forwarded-For \$proxy_add_x_forwarded_for;
        proxy_set_header X-Forwarded-Proto \$scheme;
        proxy_set_header Authorization \$http_authorization;
        proxy_http_version 1.1;
        proxy_set_header Upgrade \$http_upgrade;
        proxy_set_header Connection "upgrade";
    }

    # Drop common internet exploit scans at the edge.
    location ~* (cgi-bin|phpunit|eval-stdin\.php|think\\app|/containers/json|/\.env|/wp-admin|/wp-login|/vendor/) {
        return 444;
    }

    # HA integration API on 8124 (no redirect; preserves Authorization headers).
    location ^~ /api/ {
        proxy_pass http://127.0.0.1:8000;
        proxy_set_header Host \$host;
        proxy_set_header X-Real-IP \$remote_addr;
        proxy_set_header X-Forwarded-For \$proxy_add_x_forwarded_for;
        proxy_set_header X-Forwarded-Proto \$scheme;
        proxy_set_header Authorization \$http_authorization;
        proxy_http_version 1.1;
        proxy_set_header Upgrade \$http_upgrade;
        proxy_set_header Connection "upgrade";
    }

    # All other trusted traffic on 8124 goes to Django (single-domain compatibility mode).
    location / {
        proxy_pass http://127.0.0.1:8000;
        proxy_set_header Host \$host;
        proxy_set_header X-Real-IP \$remote_addr;
        proxy_set_header X-Forwarded-For \$proxy_add_x_forwarded_for;
        proxy_set_header X-Forwarded-Proto \$scheme;
        proxy_set_header Authorization \$http_authorization;
        proxy_http_version 1.1;
        proxy_set_header Upgrade \$http_upgrade;
        proxy_set_header Connection "upgrade";
    }
}

server {
    listen 8125 ssl;
    server_name _;

    access_log /dev/stdout ocw;
    error_log /dev/stderr info;

    ssl_certificate $CERT_FILE;
    ssl_certificate_key $KEY_FILE;

    # Only hosts listed in http_domain/tcu_domain are served.
    if (\$is_trusted_host = 0) {
        return 444;
    }

    # Serve static files directly
    location /static/ {
        alias /opt/opencarwings/staticfiles/;
        expires 30d;
        add_header Cache-Control "public, no-transform";
    }

    # Drop common internet exploit scans at the edge.
    location ~* (cgi-bin|phpunit|eval-stdin\.php|think\\app|/containers/json|/\.env|/wp-admin|/wp-login|/vendor/) {
        return 444;
    }

    # Optional direct HTTPS access to app/API.
    location / {
        proxy_pass http://127.0.0.1:8000;
        proxy_set_header Host \$host;
        proxy_set_header X-Real-IP \$remote_addr;
        proxy_set_header X-Forwarded-For \$proxy_add_x_forwarded_for;
        proxy_set_header X-Forwarded-Proto https;
        proxy_set_header X-Forwarded-Port 8125;
        proxy_set_header Authorization \$http_authorization;
        proxy_http_version 1.1;
        proxy_set_header Upgrade \$http_upgrade;
        proxy_set_header Connection "upgrade";
    }
}
EOF

# Source bashio log enhancements
source /etc/cont-init.d/00-log-fix.sh

# Generate dynamic Nginx configuration
# ... (Nginx config generation) ...
# (cat > /etc/nginx/http.d/default.conf <<EOF ... EOF already handled)

bashio::log.info "Nginx configured. Starting and monitoring logs..."

# Wait for Daphne to be ready (prevent 111: Connection refused)
bashio::log.info "Waiting for OpenCarwings backend to start..."
timeout=30
while ! nc -z 127.0.0.1 8000; do
    timeout=$((timeout - 1))
    if [ $timeout -eq 0 ]; then
        bashio::log.error "Timed out waiting for OpenCarwings backend!"
        break
    fi
    sleep 1
done
bashio::log.info "Backend is ready. Starting Nginx..."

# Start nginx and handle its output
# We pipe access/error logs through bashio log cleaners to match standard format
nginx -g "daemon off;" 2>&1 | while read -r line; do
    if [[ "$line" == *"[error]"* ]]; then
        bashio::log.error "$line" "nginx"
    elif [[ "$line" == *"[warn]"* ]]; then
        bashio::log.warning "$line" "nginx"
    else
        # Access logs and non-error nginx lines stay in debug by request.
        bashio::log.debug "$line" "nginx"
    fi
done
