# opencarwings/rootfs/etc/services.d/periodic/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash

cd /opt/opencarwings || exit 1

# Wait for Django setup to complete
# Use the dedicated health check service on port 8899
while ! curl -s http://localhost:8899/health > /dev/null; do
  bashio::log.info "Waiting for OpenCarwings service to be ready..."
  sleep 5
done

bashio::log.info "Starting OpenCarwings periodic update scheduler..."

# Set environment variables
export PSQL_DATABASE=carwings
export PSQL_USER=carwings_user
export PSQL_PASSWORD=secure_password
export PSQL_DATABASE_HOST=localhost
export PSQL_DATABASE_PORT=5432
export REDIS_HOST=localhost
export REDIS_PORT=6379
export PYTHONPATH="/opt/opencarwings${PYTHONPATH:+:$PYTHONPATH}"
export DJANGO_SETTINGS_MODULE=carwings.settings

if [ "$(bashio::config 'log_level')" = "debug" ]; then
    export DEBUG=True
else
    export DEBUG=False
fi

while true; do
  python3 manage.py tcuperiodicupdate > /dev/null 2>&1
  python3 manage.py flushexpiredtokens > /dev/null 2>&1
  bashio::log.info "Periodic background tasks executed (timeouts checked & tokens flushed)."
  sleep 60
done
