# opencarwings/rootfs/etc/services.d/healthcheck/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash

# --- Unified Logging Fix ---
# shellcheck source=/etc/cont-init.d/00-log-fix.sh
source /etc/cont-init.d/00-log-fix.sh

bashio::log.info "Starting healthcheck service..."
exec python3 /opt/opencarwings/healthcheck.py 2>&1