# opencarwings/rootfs/etc/services.d/tcuserver/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash

cd /opt/opencarwings || exit 1

bashio::log.info "Starting OpenCarwings TCU Socket Server..."

# --- Load Shared Environment ---
# shellcheck source=/etc/opencarwings/env.sh
source /etc/opencarwings/env.sh

# Run the management command
# It listens on all interfaces (0.0.0.0) within the container
exec python3 manage.py tcuserver 0.0.0.0
